# scripts/vm/alpine-tmux/init.d/tmux-tuto-serial-perms
#!/sbin/openrc-run

description="Set ttyS1/ttyS2 permissions for tmux-tuto serial channels"

depend() {
  need devfs mdev
  before agetty
}

start() {
  if [ -e /dev/ttyS1 ]; then
    chgrp tuto /dev/ttyS1 || true
    chmod 0666 /dev/ttyS1 || true
  fi
  if [ -e /dev/ttyS2 ]; then
    chgrp tuto /dev/ttyS2 || true
    chmod 0666 /dev/ttyS2 || true
  fi
}
